# supervisor-develop/supervisord-headphones.ini
[program:headphones]
command=/usr/bin/env python2 /headphones/Headphones.py --datadir=/config --nolaunch
user=headphones
autostart=true
autorestart=false
umask=002
startretries=3
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
